# openwrt/files/xfrpc.ini
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

NAME=xfrpc
PROG=/user/bin/$NAME


handle_xfrpc() {
  local name="$1"
  local config="$2"
}

service_triggers() {
        procd_add_reload_trigger "$NAME"
}

start_service() {
        local init_cfg=" "
        local conf_file="/var/etc/$NAME.ini"

        > "$conf_file"
        config_load "$NAME"

        local disabled
        uci_validate_section "$NAME" init "" 'disabled:bool:1'
        if [ $disabled = 1 ]; then
                return
        fi

        config_foreach handle_xfrpc xfrpc "$conf_file"

        procd_open_instance
        procd_set_param command "$PROG" -c "$conf_file"
        procd_set_param file "$conf_file"
        procd_set_param respawn
        procd_close_instance
}

reload_service() {
        stop
        start
}
